# config/autostart/arbtt-capture.desktop
[Desktop Entry]
Name=arbtt-capture
Comment=Automatic Rule-Based Time Tracking
Exec=arbtt-capture
Type=Application
