# config/autostart/auto-arandr.desktop
[Desktop Entry]
Name=auto-arandr
Type=Application
Exec=/home/cduquesne/.bin/auto-arandr
