# config/autostart/clipit.desktop
[Desktop Entry]
Name=clipit
Comment=clipboard manager
Exec=clipit
Type=Application
