# config/autostart/firefox.desktop
[Desktop Entry]
Type=Application
Name=firefox
Exec=firefox
Comment=start firefox
