# config/autostart/linphone.desktop
[Desktop Entry]
Type=Application
Comment=Gnome rozhraní pro linphone, internetový telefon podporující SIP.
Exec=/usr/bin/linphone
Name=linphone
