# config/autostart/nm-applet.desktop
[Desktop Entry]
Comment=network monitor and control GUI applet
Exec=/usr/bin/nm-applet
Type=Application
Name=nm-applet
