# config/autostart/pasystray.desktop
[Desktop Entry]
Exec=/usr/bin/pasystray
Comment=PulseAudio controller for the system tray
Type=Application
Name=pasystray
