# config/autostart/pidgin.desktop
[Desktop Entry]
Exec=/usr/bin/pidgin
Name=pidgin
Type=Application
Comment=Perl extension for the Pidgin instant messenger.
