# config/autostart/qrclip.desktop
[Desktop Entry]
Type=Application
Name=qrclip
Exec=qrclip
Comment=start qrclip
