# config/autostart/setxkbmap.desktop
[Desktop Entry]
Type=Application
Name=setxkbmap
Exec=setxkbmap us -variant intl
Comment=set the keyboard using the X Keyboard Extension
