# config/autostart/shutter.desktop
[Desktop Entry]
Name=shutter
Comment=Feature-rich Screenshot Tool
Type=Application
Exec=/usr/bin/shutter --min_at_startup
