# config/autostart/syndaemon.desktop
[Desktop Entry]
Comment=Deactivate the keyboard while typing
Name=syndaemon
Exec=syndaemon
Type=Application
