# config/autostart/syndaemon.desktop
[Desktop Entry]
Comment=Deactivate the keyboard while typing
Name=syndaemon
Exec=syndaemon -i 0.5 -t -K -R
Type=Application
