# config/autostart/unclutter.desktop
[Desktop Entry]
Comment=hide the mouse
Name=unclutter
Exec=unclutter
Type=Application
