# config/autostart/update-notifier.desktop
[Desktop Entry]
Name=update-notifier
Type=Application
Exec=/usr/bin/update-notifier
