# config/autostart/xfce4-panel.desktop
[Desktop Entry]
Name=xfce4-panel
Comment=xfce4 panel
Exec=/usr/bin/xfce4-panel --disable-wm-check
Type=Application
