# config/autostart/xss-lock.desktop
[Desktop Entry]
Exec=/usr/bin/xss-lock i3lock
Name=xss-lock
Type=Application
Comment=fire up programs in case of user inactivity under X
