# config/autostart/xss-lock.desktop
[Desktop Entry]
Exec=xss-lock -- i3lock-pixelate -n -f
Name=xss-lock
Type=Application
Comment=fire up programs in case of user inactivity under X
