# config/autostart/xss-lock.desktop
[Desktop Entry]
Exec=xss-lock -l -- i3lock-blur
Name=xss-lock
Type=Application
Comment=fire up programs in case of user inactivity under X
